# dot_config/systemd/user/email-notify@.service
[Unit]
Description=Email notification for %i

[Service]
Type=oneshot
ExecStart=%h/.local/bin/notify-systemd-failure %i
